# certbot.service
[Unit]
Description=Certbot
[Service]
Type=oneshot
ExecStart=/home/pi/oprint/bin/certbot -q renew
PrivateTmp=true
