# certbot.service
[Unit]
Description=Certbot
[Service]
Type=oneshot
ExecStart=/usr/local/bin/certbot -q renew
PrivateTmp=true
